# device/rockchip/common/overlays/rootfs/usb-gadget/lib/udev/rules.d/61-usb-gadget.rules
SUBSYSTEM=="android_usb", ACTION=="change", RUN+="/usr/bin/usb-gadget update"
SUBSYSTEM=="udc", RUN+="/usr/bin/usb-gadget update"
