# device/rockchip/common/overlays/x21-misc/etc/udev/rules.d/99-sdcard.rules
ACTION=="add", KERNEL=="mmcblk1p1", SUBSYSTEM=="block", RUN+="/usr/bin/update-from-sd.sh"